# lib/file-templates/supervisor.conf
[program:{{name}}]
command=/home/services/{{name}}.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/{{name}}.log
stderr_logfile=/var/log/{{name}}.log

stopsignal=KILL
stopasgroup=true